# autossh/rootfs/etc/services.d/autossh/finish
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# s6-overlay finish script: runs after the service exits
# $1 = exit code (256 if killed by signal)

bashio::log.warning "Service exited with code ${1}. Restarting in 60 seconds..."
sleep 60
